# liquidctl/liquidctl-info.timer
[Unit]
Description=Timer for liquidctl-info service

[Timer]
OnUnitInactiveSec=2s
OnBootSec=10s

[Install]
WantedBy=timers.target
